# supervisor.conf
[program:fflow]
command=/var/www/fflow/localpy/bin/python /var/www/fflow/fflow.py
user=www-data
autorestart=true
startretries=1000000
